# ha-pi.service
[Unit]
Description=Home Assistant Light Control Touchscreen UI
After=network.target

[Service]
Type=simple
ExecStartPre=/bin/sh -c 'killall fbcp 2>/dev/null || true'
ExecStartPre=/bin/sh -c 'for vt in /sys/class/vtconsole/vtcon*/bind; do echo 0 > "$vt" 2>/dev/null; done || true'
ExecStartPre=/bin/sh -c 'echo 0 > /sys/class/graphics/fbcon/cursor_blink 2>/dev/null || true'
ExecStartPre=/bin/sh -c 'setterm --blank force --cursor off > /dev/tty1 2>/dev/null || true'
ExecStart=REPLACE_HOME/ha-pi/ha_lights /etc/ha_lights.conf
Restart=on-failure
RestartSec=5
User=REPLACE_USER
Group=REPLACE_USER

# Hardware access
SupplementaryGroups=spi gpio video

[Install]
WantedBy=multi-user.target
